# deploy/launchd/com.chatbi.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chatbi.worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yangjinlong/app/PythonProject/ChatBiTest/.venv/bin/python</string>
        <string>/Users/yangjinlong/app/PythonProject/ChatBiTest/worker.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yangjinlong/app/PythonProject/ChatBiTest</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/yangjinlong/app/PythonProject/ChatBiTest/logs/launchd-worker.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yangjinlong/app/PythonProject/ChatBiTest/logs/launchd-worker.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
